# services/backupbot/run
#!/usr/bin/env bash
set -e

LOG_FILE="/config/log/backupbot.log"
mkdir -p "$(dirname "$LOG_FILE")"

echo "[BACKUPBOT_INFO] Starting PostgreSQL backup loop service..." | tee -a "$LOG_FILE"

INTERVAL_HOURS="${INTERVAL_HOURS:-24}"

# Stay alive forever
while true; do
  echo "[BACKUPBOT_INFO] Triggering backup.sh at $(date)" | tee -a "$LOG_FILE"
  /usr/local/bin/backup.sh "$LOG_FILE"
  echo "[BACKUPBOT_INFO] Sleeping for $INTERVAL_HOURS hours..." | tee -a "$LOG_FILE"
  sleep "${INTERVAL_HOURS}h"
done
