# services/backupbot/run
#!/usr/bin/with-contenv bash
set -e

echo "[BACKUPBOT_INFO] Starting PostgreSQL backup loop service..."

INTERVAL_HOURS="${INTERVAL_HOURS:-24}"
STATE_FILE="/config/last_backup_date"
LOG_FILE="/config/log/pgbackup.log"

mkdir -p "$(dirname "$STATE_FILE")" "$(dirname "$LOG_FILE")"

while true; do
  TODAY=$(date +%F)

  # Check if a backup already ran today
  if [[ -f "$STATE_FILE" && "$(cat "$STATE_FILE")" == "$TODAY" ]]; then
    echo "[BACKUPBOT_INFO] Backup already completed today ($TODAY). Skipping."
  else
    echo "[BACKUPBOT_INFO] Triggering backup.sh at $(date)"
    /usr/local/bin/backup.sh "$LOG_FILE"
    echo "$TODAY" >"$STATE_FILE"
    echo "[BACKUPBOT_INFO] Backup completed and date recorded."
  fi

  echo "[BACKUPBOT_INFO] Sleeping for $INTERVAL_HOURS hours..."
  sleep "${INTERVAL_HOURS}h"
done
